# docker/coturn/rootfs/etc/systemd/system/coturn.service.d/10_service_execstartpre.conf
[Service]
ExecStartPre=/usr/local/sbin/coturn.sh
