# docker/worker/kolab-worker.service
[Unit]
Description=Kolab Worker Service

[Service]
Type=simple
ExecStart=/usr/local/sbin/kolab-worker.sh

[Install]
WantedBy=multi-user.target
